# share/systemd/coding-plans-today.timer
[Unit]
Description=coding-plans-waybar: run today's Claude backfill every 5 min

[Timer]
# OnActiveSec: fire 30 s after the timer itself is activated (i.e. right
#   after `systemctl enable --now`) — otherwise a fresh install waits up
#   to 5 min before the first backfill and Today reads 0.
# OnUnitActiveSec: then repeat every 5 min based on the service's last run.
# OnBootSec: also cover reboots (user session start).
OnActiveSec=30s
OnBootSec=1min
OnUnitActiveSec=5min
Unit=coding-plans-today.service
Persistent=true

[Install]
WantedBy=timers.target
